# docker-containers/microservice_flask/supervisor/run_app.conf
[program:app]
command=python3 /opt/microservice_flask/run_app.py
