# microservice_templates/microservice_python3_template/supervisor/_MICROSERVICE_PYTHON3_TEMPLATE_.conf
[program:_MICROSERVICE_PYTHON3_TEMPLATE_]
directory=/opt/_MICROSERVICE_PYTHON3_TEMPLATE_/src
command=python3 _MICROSERVICE_PYTHON3_TEMPLATE_.py
