# microservice_templates/microservice_dotnet_template/supervisor/_MICROSERVICE_DOTNET_TEMPLATE_.conf
[program:_MICROSERVICE_DOTNET_TEMPLATE_]
directory=/opt/_MICROSERVICE_DOTNET_TEMPLATE_
command=sh -c 'dotnet restore && exec dotnet run'
