# docker-containers/microservice_php_jammy/supervisor/cron.conf
[program:cron]
command=/usr/sbin/cron -f
autorestart=true
