# microservice_templates/microservice_python_template/supervisor/_MICROSERVICE_PYTHON_TEMPLATE_.conf
[program:_MICROSERVICE_PYTHON_TEMPLATE_]
directory=/opt/_MICROSERVICE_PYTHON_TEMPLATE_/src
command=python _MICROSERVICE_PYTHON_TEMPLATE_.py 80
