# microservice_templates/microservice_node_template/supervisor/_MICROSERVICE_NODE_TEMPLATE_.conf
[program:_MICROSERVICE_NODE_TEMPLATE_]
directory=/opt/_MICROSERVICE_NODE_TEMPLATE_/src
command=nodejs _MICROSERVICE_NODE_TEMPLATE_.js
