# docker-containers/microservice/supervisor/local_magellan.conf
[program:local_magellan]
directory=/opt/microservice/src/local_magellan
command=python local_magellan.py
