# docker-containers/microservice/supervisor/register_in_service_discovery.conf
[program:register_in_service_discovery]
command=microservice register 80
startsecs=0
priority=10
