# docker-containers/microservice/supervisor/require.conf
[program:require]
directory=/opt/microservice/src/local_magellan
command=python require_service.py -c service_discovery.json
startsecs=0
