# docker-containers/microservice/supervisor/save_environment_variables.conf
[program:save_environment_variables]
directory=/opt/microservice/src
command=python save_environment_variables.py
autorestart=false
startsecs=0
