# docker-containers/microservice_jammy/packaging/microservice/etc/supervisor/conf.d/require.conf
[program:require]
command=microservice require -c service_discovery.json
startsecs=0
