# install/openrc_armada
#!/sbin/runscript

name='armada'

depend() {
        after docker
        need docker
        use docker
}

start() {
        /usr/local/bin/armada-runner start
}

stop() {
        /usr/local/bin/armada-runner stop
}
