# docker-containers/microservice/packaging/microservice/etc/supervisor/conf.d/armada_agent.conf
[program:armada_agent]
command=microservice agent
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
