# docker-containers/microservice/packaging/microservice/etc/supervisor/conf.d/local_magellan.conf
[program:local_magellan]
command=microservice local-magellan
exitcodes=0
