# docker-containers/microservice_php/supervisor/apache2.conf
[program:apache2]
command=/opt/microservice_php/run_apache2.sh
