# armada_backend/supervisor/armada.conf
[program:armada]
directory=/opt/armada-docker/armada_backend
command=bash share_armada.sh
autorestart=false
startsecs=0
