# armada_backend/supervisor/backup_running_containers.conf
[program:backup_running_containers]
directory=/opt/armada-docker/armada_backend
command=bash backup_running_containers.sh
