# armada_backend/supervisor/cleaner.conf
[program:cleaner]
directory=/opt/armada-docker/armada_backend
command=python cleaner.py
