# armada_backend/supervisor/consul.conf
[program:consul]
directory=/opt/armada-docker/armada_backend
command=python3 run_consul.py
autorestart=true
