# armada_backend/supervisor/docker_client.conf
[program:install_docker]
directory=/opt/armada-docker/armada_backend
command=python install_docker_client.py
autorestart=false
