# armada_backend/supervisor/hermes-init.conf
[program:hermes_init]
directory=/opt/armada-docker/armada_backend
command=python3 hermes_init.py
autorestart=false
startsecs=0
