# armada_backend/supervisor/kv_store_updater.conf
[program:kv_store_updater]
directory=/opt/armada-docker/armada_backend
command=python3 kv_store_updater.py
