# armada_backend/supervisor/load_running_services.conf
[program:load_running_services]
directory=/opt/armada-docker/armada_backend
command=python load_running_services.py /opt/armada/running_containers_parameters.json
autorestart=false
startsecs=0